# init/service/rc.local
#!/bin/bash
#sysctl -p
#/usr/local/bin/ss-go start &
#/usr/local/bin/ss-go-nofw start &
